# libvirtd/supervisor.conf
[supervisord]

[program:sshd]
command=/usr/sbin/sshd -D -p 2222

[program:libvirtd]
command=/usr/sbin/libvirtd

